# supervisord.conf
; supervisor config file

[inet_http_server]
port=127.0.0.1:10000

[supervisord]
logfile = /home/team41/.supervisord.log

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=http://127.0.0.1:10000 ; use a unix:// URL  for a unix socket

; The [include] section can just contain the "files" setting.  This
; setting can list multiple files (separated by whitespace or
; newlines).  It can also contain wildcards.  The filenames are
; interpreted as relative to this file.  Included files *cannot*
; include files themselves.

[program:mlacademy]
directory = /home/team41/production_django/
command = authbind gunicorn --worker-class=eventlet --workers=2 --name=backend --timeout=60 --bind='0.0.0.0:80' mlacademy.wsgi
user = team41
stdout_logfile = /backend.log
redirect_stderr = true